# certs-update.service
[Unit]
Description=Update letsencrypt certificates
[Service]
Type=oneshot
ExecStart=cd /root/letsencrypt-power-update && ./updateAll.sh
